# ping_probe_client.cron.d
# PingProbe runner

* * * * * root /opt/ping_probe/run
