# cron.d/ping_probe
# PingProbe runner

* * * * * pinger /home/pinger/ping_client/run
